# nordvpn-client/etc/service/nordvpn/run
#!/bin/sh
set -e

echo "Starting nordvpnd in the foreground..."
# By default, nordvpnd might log to /var/log/nordvpnd.log, but let's log to stdout:
exec /usr/sbin/nordvpnd --log /dev/stdout
