# nordvpn-client/etc/service/nordvpn_init/run
#!/bin/sh

# Wait for the nordvpnd socket, which indicates the daemon is ready
while [ ! -S /run/nordvpn/nordvpnd.sock ]; do
    echo "Waiting for NordVPN daemon to create /run/nordvpn/nordvpnd.sock..."
    sleep 2
done

echo "NordVPN daemon socket found, proceeding with login/connect..."

# Example token-based login (or skip if you already have credentials saved):
if [ -n "${NORDVPN_TOKEN}" ]; then
  echo "Logging in via token..."
  nordvpn login --token "${NORDVPN_TOKEN}" || true
fi

# Example settings
if [ -n "${TECHNOLOGY}" ]; then
  nordvpn set technology "${TECHNOLOGY}"
fi

# Whitelist port 22 (SSH):
#nordvpn whitelist add port 22

# Or whitelist your LAN subnet if you want full LAN access:
# Suppose your LAN is 192.168.0.x
nordvpn whitelist add subnet 192.168.0.0/16


# Some users like killswitch, autoconnect, etc.
nordvpn set killswitch on
nordvpn set autoconnect on

# Connect (no argument = best server, or specify e.g. "Denmark", "Germany", etc.)
nordvpn connect

# Keep this service from exiting, so runit doesn't constantly restart it
echo "NordVPN is connected. Service will now idle..."
exec tail -f /dev/null
